# crontab
*/3 * * * * /usr/local/cnupmsave/upload.sh